# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:xvfb]
command=Xvfb :99 -screen 0 1280x720x24
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:chromium]
command=chromium --no-sandbox --disable-gpu --disable-software-rasterizer --display=:99 --user-data-dir=/data/chromium-profile "https://accounts.google.com/ServiceLogin?service=youtube"
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
startsecs=5

[program:x11vnc]
command=/app/start-vnc.sh
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:websockify]
command=websockify --web /usr/share/novnc 6080 localhost:5900
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:node]
command=node /app/dist/server.js
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
